# supervisor/laravel-reverb-c3ers.conf
; Laravel Reverb — development (501c3ers.com)
; https://laravel.com/docs/13.x/reverb#process-management
;
; Install (root):
;   cp supervisor/laravel-reverb-c3ers.conf /etc/supervisord.d/laravel-reverb-c3ers.ini
;   supervisorctl reread && supervisorctl update && supervisorctl start laravel-reverb-c3ers:*
;
; Server .env (see .env.example):
;   REVERB_SERVER_HOST=127.0.0.1   — bind address for reverb:start
;   REVERB_SERVER_PORT=8080
;   REVERB_HOST=127.0.0.1            — PHP broadcasts (internal)
;   REVERB_PORT=8080
;   VITE_REVERB_*                    — browser wss via nginx :443
;
; Optional: set minfds=10000 in /etc/supervisord.conf [supervisord] section.

[program:laravel-reverb-c3ers]
process_name=%(program_name)s
command=/opt/cpanel/ea-php84/root/usr/bin/php /home/c3ers/public_html/artisan reverb:start
directory=/home/c3ers/public_html
user=c3ers
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/home/c3ers/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
