# supervisor/laravel-reverb-believeinunity.conf
; Laravel Reverb — production (believeinunity.org)
; https://laravel.com/docs/13.x/reverb#process-management
;
; Install (root):
;   cp supervisor/laravel-reverb-believeinunity.conf /etc/supervisord.d/laravel-reverb-believeinunity.ini
;   supervisorctl reread && supervisorctl update && supervisorctl start laravel-reverb-believeinunity:*

[program:laravel-reverb-believeinunity]
process_name=%(program_name)s
command=/opt/cpanel/ea-php84/root/usr/bin/php /home/believeinunity/public_html/artisan reverb:start
directory=/home/believeinunity/public_html
user=believeinunity
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/home/believeinunity/logs/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
